# operational/systemd/rubin-node.service
[Unit]
Description=RUBIN node with verified wolfcrypt shim (local single-host)
After=network.target

[Service]
Type=simple
#
# Template for dev/test single-host use. Adjust paths for your environment.
#
WorkingDirectory=/opt/rubin/rubin-protocol
EnvironmentFile=/etc/rubin/rubin-node.env
ExecStart=/bin/bash -lc '\
  set -euo pipefail && \
  source <(./crypto/wolfcrypt/fetch_shim_from_ghcr.sh "${GHCR_REF}" /var/lib/rubin/wolfcrypt-shim --export-env) && \
  if [ "${CLIENT}" = "go" ]; then \
    go run ./clients/go/node chain-id; \
  else \
    cargo run -p rubin-node -- chain-id; \
  fi'
Restart=on-failure

[Install]
WantedBy=multi-user.target
